# 2017-05-10 Python/Novel&Image/cl/spider.cron
0 */1 * * * cd /noval/ && python3 spider.py > noval.log
